# agent-server/launchd/com.aios.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aios.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string><PROJECT_ROOT>/agent-server/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string><PROJECT_ROOT>/agent-server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>America/Mexico_City</string>
        <key>HOME</key>
        <string><HOME></string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>StandardOutPath</key>
    <string><HOME>/Library/Logs/aios/agent-server.log</string>

    <key>StandardErrorPath</key>
    <string><HOME>/Library/Logs/aios/agent-server.err.log</string>
</dict>
</plist>
